# alloy.service
[Unit]
Description=Grafana Alloy
After=network.target
Wants=network.target

[Service]
Type=simple
WorkingDirectory=%h
ExecStart=%h/.local/bin/alloy run %h/registry/alloy/config.alloy --server.http.listen-addr=0.0.0.0:12345 --storage.path=%h/.local/share/alloy
Restart=always
RestartSec=5

[Install]
WantedBy=default.target