# docker-proxy.service
[Unit]
Description=Docker Socket Proxy for Rootless Docker
After=docker.service
Requires=docker.service

[Service]
Type=simple
User=madmin
Group=madmin
ExecStart=/usr/bin/socat TCP-LISTEN:2375,bind=127.0.0.1,reuseaddr,fork UNIX-CONNECT:/run/user/1000/docker.sock
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target